# config/supervisord.conf
[supervisord]
nodaemon=true

[program:vpnserver]
command=bash -c 'cp /opt/vpnserver/* /mnt/softether/; exec /mnt/softether/vpnserver execsvc'
